# resources/systemd/module-gateway.service.in
[Unit]
Description=Fleet protocol's Module Gateway
After=network.target nss-lookup.target

[Service]
Type=simple
Restart=always
RestartSec=20
WorkingDirectory=@CMAKE_INSTALL_SYSCONFDIR@
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/module-gateway-app --config-path=@CMAKE_INSTALL_SYSCONFDIR@/bringauto/module-gateway/config.json

[Install]
WantedBy=multi-user.target
